# scripts/podman/genoma.container.in
# Genome gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into genome's ~/.config/containers/systemd/
# {{GENOME_HOME}} is replaced at install time.

[Unit]
Description=Genome gateway (rootless Podman)

[Container]
Image=genome:local
ContainerName=genome
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{GENOME_HOME}}/.genome:/home/node/.genome
EnvironmentFile={{GENOME_HOME}}/.genome/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
